# Pi/systemd/eink-ip.timer
# Re-run eink-ip.service periodically so the displayed IP follows DHCP
# changes without needing a reboot.  Eink refresh takes ~3 s and the
# panel handles tens of thousands of full refreshes — 5-minute polling
# is comfortable.
#
# A change-only refresh would be nicer; for the MVP we just re-render
# unconditionally.  Swap to a NetworkManager dispatcher hook later if
# you want event-driven updates.

[Unit]
Description=Refresh LaserHAT eink IP display every 5 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=5min
Unit=eink-ip.service

[Install]
WantedBy=timers.target
