# Pi/systemd/eink-gui.service
# LaserHAT eink GUI — long-running daemon.
#
#   sudo cp Pi/systemd/eink-gui.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now eink-gui.service
#
# A broker client: connects to laserhat-broker.service over the Unix socket
# and doesn't open the UART, so it runs alongside the web GUI. Edit User=,
# Group=, WorkingDirectory=, and the ExecStart= python prefix below if your
# Pi account isn't 'kemerelab' or the venv lives elsewhere.

[Unit]
Description=LaserHAT eink GUI (broker client + display)
After=network-online.target laserhat-broker.service
Wants=network-online.target
Requires=laserhat-broker.service

[Service]
Type=simple
User=kemerelab
Group=kemerelab
WorkingDirectory=/home/kemerelab/Code/LaserDriver/Pi
Environment=LASERHAT_SOCK=/run/laserhat/broker.sock
ExecStart=/home/kemerelab/.venvs/laserhat/bin/python eink_gui.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
