# Pi/systemd/laserhat-web.service
# LaserHAT web GUI (Flask) — long-running daemon.
#
# A broker client: it connects to laserhat-broker.service over the Unix
# socket, so it no longer needs exclusive UART access and can run at the
# same time as eink-gui.service.
#
#   sudo cp Pi/systemd/laserhat-web.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now laserhat-web.service
#
# Then browse to http://<pi>:8080/ from any device on the LAN.
#
# Edit User=, Group=, paths and ExecStart= if your Pi account isn't
# 'kemerelab' or the venv lives somewhere other than ~/.venvs/laserhat.

[Unit]
Description=LaserHAT web GUI (Flask)
After=network-online.target laserhat-broker.service
Wants=network-online.target
Requires=laserhat-broker.service

[Service]
Type=simple
User=kemerelab
Group=kemerelab
WorkingDirectory=/home/kemerelab/Code/LaserDriver/Pi
Environment=PORT=8080
Environment=LASERHAT_SOCK=/run/laserhat/broker.sock
ExecStart=/home/kemerelab/.venvs/laserhat/bin/python web_app.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
