# Pi/systemd/oled-gui.service
# LaserHAT OLED GUI — long-running daemon (Adafruit 2.23" 128x32 bonnet).
#
#   sudo cp Pi/systemd/oled-gui.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now oled-gui.service
#
# A broker client: connects to laserhat-broker.service over the Unix socket
# and doesn't open the UART, so it runs alongside the web GUI. Edit User=,
# Group=, WorkingDirectory=, and the ExecStart= python prefix below if your
# Pi account isn't 'kemerelab' or the venv lives elsewhere.

[Unit]
Description=LaserHAT OLED GUI (broker client + display)
After=network-online.target laserhat-broker.service
Wants=network-online.target
Requires=laserhat-broker.service

[Service]
Type=simple
User=kemerelab
Group=kemerelab
WorkingDirectory=/home/kemerelab/Code/LaserDriver/Pi
Environment=LASERHAT_SOCK=/run/laserhat/broker.sock
ExecStart=/home/kemerelab/.venvs/laserhat/bin/python oled_gui.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
